# automation/com.vincemaguire.podcast-feed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vincemaguire.podcast-feed</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/vincemaguire/My Drive/Podcast Automation/watch-and-deploy.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/vincemaguire/My Drive/Podcast Automation</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/vincemaguire/My Drive/Podcast Automation/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/vincemaguire/My Drive/Podcast Automation/logs/launchd.err.log</string>
</dict>
</plist>
